# ops/systemd/detecdiv-worker.service
[Unit]
Description=DetecDiv Hub Worker
After=network.target

[Service]
Type=simple
User=detecdiv
WorkingDirectory=/srv/detecdiv/detecdiv-hub
Environment=DETECDIV_HUB_DATABASE_URL=postgresql+psycopg://detecdiv:detecdiv@localhost:5432/detecdiv_hub
ExecStart=/srv/detecdiv/detecdiv-hub/.venv/bin/python worker/run_worker.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

