# ops/systemd/detecdiv-api.service
[Unit]
Description=DetecDiv Hub API
After=network.target

[Service]
Type=simple
User=detecdiv
WorkingDirectory=/srv/detecdiv/detecdiv-hub
EnvironmentFile=/etc/detecdiv-hub/detecdiv-hub.env
ExecStart=/srv/detecdiv/detecdiv-hub/.venv/bin/uvicorn api.app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
